# deploy/observability/systemd/alloy.service.d/override.conf
[Service]
User=alloy
Group=alloy
UMask=0027

# PID 1 reads this root-only file before Alloy drops to its dedicated user.
EnvironmentFile=/etc/alloy/alloy.env
Environment="CUSTOM_ARGS=--server.http.listen-addr=127.0.0.1:12345 --disable-reporting"

# Keep the UI, component graph, health, and /metrics local to the VM. The
# package's --storage.path=/var/lib/alloy/data persists file positions and the
# Prometheus remote-write WAL across restarts.
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

ReadOnlyPaths=/etc/alloy
ReadOnlyPaths=/var/log/mapleland-api
ReadWritePaths=/var/lib/alloy
